# services/mosquitto/finish
#!/bin/bash
echo 'mosquitto exited... restarting'
