# services/mosquitto/finish
#!/bin/bash
echo 'icecast exited... restarting'
