# services/mosquitto/run
#!/usr/bin/with-contenv /bin/bash
echo 'Starting mosquitto'
exec mosquitto -c /mosquitto/mosquitto.conf

