# services/cron/run
#!/usr/bin/with-contenv /bin/bash
exec cron -f
